# airootfs/etc/systemd/system/xkcd.service
[Unit]
Description=Download last xkcd picture

[Service]
Type=simple
User=mc
ExecStart=/home/mc/sh_script/get_image.sh
